# ansible/templates/strapi/strapi.service.j2
[Unit]
Description=Strapi App Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Environment variables
EnvironmentFile={{ ENV_FILE_PATH }}
Environment=NODE_ENV=production

# App directory (contains package.json, node_modules, build/)
WorkingDirectory={{ service_working_dir }}

# Start Strapi
ExecStart=/usr/bin/node {{ service_working_dir }}/node_modules/@strapi/strapi/bin/strapi.js start

# Restart policy
Restart=always
RestartSec=5

# File descriptor limit (useful for DB connections)
LimitNOFILE=65535

# Run as non-root
User={{ app_user }}
Group={{ app_group }}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
